# deploy/systemd/gsv-tracker.service
# GSV Tracker daily collection service (user unit).
# Install: see deploy/README.md
[Unit]
Description=GSV Tracker: collect today's due cities
Wants=network-online.target
After=network-online.target
# Opt-in raw safety net: email whenever this unit exits nonzero. run-due
# returns nonzero on ANY failed city, so the scheduler's own threshold email
# ([alerts] in scheduler.toml) is usually the better signal; enable this too
# only if you want a belt-and-suspenders notice on every failure/crash.
#OnFailure=gsv-tracker-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=%h/gsv-tracker
# GMAPS_API_KEY lives here (same format as the repo .env: KEY=value)
EnvironmentFile=%h/gsv-tracker/.env
ExecStart=%h/gsv-tracker/.venv/bin/python -m gsv_metadata_tracker.scheduler run-due --config %h/gsv-tracker/config/scheduler.toml
# A full day of cities can legitimately take hours
TimeoutStartSec=12h
Nice=10
