# deploy/systemd/streetscape-tracker.service
# Streetscape Tracker daily collection service (user unit).
# Install: see deploy/README.md
[Unit]
Description=Streetscape Tracker: collect today's due cities
Wants=network-online.target
After=network-online.target
# Opt-in raw safety net: email whenever this unit exits nonzero. run-due
# returns nonzero on ANY failed city, so the scheduler's own threshold email
# ([alerts] in scheduler.toml) is usually the better signal; enable this too
# only if you want a belt-and-suspenders notice on every failure/crash.
#OnFailure=streetscape-tracker-notify@%n.service

[Service]
Type=oneshot
# %h/streetscape-tracker is a symlink to the real checkout on lab storage
# (/projects/makeabilitylab/streetscape-tracker); see deploy/README.md.
WorkingDirectory=%h/streetscape-tracker
# GMAPS_API_KEY / MAPILLARY_ACCESS_TOKEN live here (repo .env format: KEY=value)
EnvironmentFile=%h/streetscape-tracker/.env
# Publish straight to the NFS-mounted docroot (no SSH) — see sync_data_to_server.sh
Environment=STREETSCAPE_PUBLISH_LOCAL=1
# NB: --config is a global arg and MUST precede the subcommand (argparse).
ExecStart=%h/streetscape-tracker/.venv/bin/python -m streetscape_metadata_tracker.scheduler --config %h/streetscape-tracker/config/scheduler.makelab1.toml run-due
# A full day of cities can legitimately take hours
TimeoutStartSec=12h

# ── Sandboxing ─────────────────────────────────────────────────────────
# The job needs exactly: the checkout (data/, logs/, the SQLite catalog)
# and the NFS-mounted public docroot. Everything else is read-only.
# PrivateUsers= is required first: the mount-namespace options below are
# unavailable to user units without it (they fail with EPERM).
# Verify after the first sandboxed night:
#   journalctl --user -u streetscape-tracker.service | grep -i -m1 'permission\|read-only'
#
# RestrictSUIDSGID= must stay OFF: the publish step runs rsync with
# --chmod=D2755 (setgid dirs, so docroot files inherit the web group), and
# RestrictSUIDSGID EPERMs any chmod that sets the SGID bit — files still
# transfer but rsync exits 23, so every nightly publish reports failure
# (observed on the first live night, 2026-07-16).
#
# Known limitation: local sendmail submission is setgid and thus
# incompatible with NoNewPrivileges (which several options here imply
# anyway), so the [alerts] 'mail' transport cannot work under this
# sandbox — s-nail must be configured to speak SMTP directly instead.
PrivateUsers=true
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Real paths, not the %h symlink — mount namespacing resolves symlinks.
ReadWritePaths=/projects/makeabilitylab/streetscape-tracker
ReadWritePaths=/cse/web/research/makelab/public/streetscape-tracker

# ── Resource guardrails ────────────────────────────────────────────────
# This host is shared with other lab services (and storage is on a shared NFS
# array). These caps keep nightly collection from starving co-tenants.
# The workload is I/O-bound async HTTP + pandas, so real peaks are modest;
# validate after the first live night and tune:
#   systemctl --user show streetscape-tracker.service -p MemoryPeak -p CPUUsageNSec
Nice=10
MemoryAccounting=true
MemoryHigh=4G
MemoryMax=8G
CPUAccounting=true
CPUQuota=400%
# Note: data lives on NFS, so IO is network, not block-device —
# IOWeight would not govern it; CPU/Memory caps + Nice are the real levers.
