# deploy/systemd/streetscape-tracker.service
# Streetscape Tracker daily collection service (user unit).
# Install: see deploy/README.md
[Unit]
Description=Streetscape Tracker: collect today's due cities
Wants=network-online.target
After=network-online.target
# Opt-in raw safety net: email whenever this unit exits nonzero. run-due
# returns nonzero on ANY failed city, so the scheduler's own threshold email
# ([alerts] in scheduler.toml) is usually the better signal; enable this too
# only if you want a belt-and-suspenders notice on every failure/crash.
#OnFailure=streetscape-tracker-notify@%n.service

[Service]
Type=oneshot
# %h/streetscape-tracker is a symlink to the real checkout on lab storage
# (/projects/makeabilitylab/streetscape-tracker); see deploy/README.md.
WorkingDirectory=%h/streetscape-tracker
# GMAPS_API_KEY / MAPILLARY_ACCESS_TOKEN live here (repo .env format: KEY=value)
EnvironmentFile=%h/streetscape-tracker/.env
# Publish straight to the NFS-mounted docroot (no SSH) — see sync_data_to_server.sh
Environment=STREETSCAPE_PUBLISH_LOCAL=1
# NB: --config is a global arg and MUST precede the subcommand (argparse).
ExecStart=%h/streetscape-tracker/.venv/bin/python -m streetscape_metadata_tracker.scheduler --config %h/streetscape-tracker/config/scheduler.makelab1.toml run-due
# A full day of cities can legitimately take hours
TimeoutStartSec=12h

# ── Resource guardrails ────────────────────────────────────────────────
# makelab1 is shared with Project Sidewalk (and its data lives on the same
# makelab2 array). These caps keep nightly collection from starving co-tenants.
# The workload is I/O-bound async HTTP + pandas, so real peaks are modest;
# validate after the first live night and tune:
#   systemctl --user show streetscape-tracker.service -p MemoryPeak -p CPUUsageNSec
Nice=10
MemoryAccounting=true
MemoryHigh=4G
MemoryMax=8G
CPUAccounting=true
CPUQuota=400%
# Note: data lives on NFS (makelab2), so IO is network, not block-device —
# IOWeight would not govern it; CPU/Memory caps + Nice are the real levers.
